# resources/gonarch-check.service
[Unit]
Description=Gonarch Check module
# Documentation=""

[Service]
Type=simple
User=root
Group=root
TimeoutStartSec=0
Restart=always
RestartSec=30s
SyslogIdentifier=gonarch-check
WorkingDirectory=/opt/gonarch
ExecStart=/opt/gonarch/check

[Install]
WantedBy = multi-user.target
